# supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock ; path to your socket file
chmod=0700
 
[supervisord]
logfile=/var/log/supervisord.log ; supervisord log file
logfile_maxbytes=50MB ; maximum size of logfile before rotation
logfile_backups=10 ; number of backed up logfiles
loglevel=error ; info, debug, warn, trace
pidfile=/var/run/supervisord.pid ; pidfile location
nodaemon=false ; run supervisord as a daemon
minfds=1024 ; number of startup file descriptors
minprocs=200 ; number of process descriptors
user=root ; default user
childlogdir=/var/log ; where child log files will live
 
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface
 
[supervisorctl]
serverurl=unix:///var/run/supervisor.sock ; use a unix:// URL for a unix socket

[program:celerypurge]
command=celery purge --config=task_config
autostart=false
autorestart=false
startsecs=0
stopwaitsecs=1

[program:celeryd]
environment=C_FORCE_ROOT="true"
command=celery worker -l INFO -B --config=task_config
stdout_logfile=/var/log/celeryd.log
stderr_logfile=/var/log/celeryd.log
autostart=false
autorestart=true
startsecs=10
stopwaitsecs=600

[program:tasks]
command=start_tasks.sh
autostart=false
autorestart=false
startsecs=0
stopwaitsecs=1

[program:monitor]
command=celery flower --address=0.0.0.0 --port=5555 --broker="amqp://celery:pass@localhost:5672/task"
stdout_logfile=/var/log/celeryflower.log
stderr_logfile=/var/log/celeryflower.log
autostart=false
autorestart=true
startsecs=10
stopwaitsecs=600

[group:celery]
programs=celerypurge,celeryd,monitor,tasks
autostart=true

[program:rest_api]
command=python3 rest_api.py --port=3010 --celery_broker="amqp://celery:pass@localhost:5672/task" --celery_backend="db+postgresql://heck:pass@localhost/heck" --celery_tables='{"task": "celery_task", "group": "celery_group"}'
stdout_logfile=/var/log/heck_rest_api.log
stderr_logfile=/var/log/heck_rest_api.log
autostart=false
autorestart=true
startsecs=10
stopwaitsecs=600

[program:web]
command=node app -b http://localhost:3010 -p 3000
directory=web
stdout_logfile=/var/log/heck_web.log
stderr_logfile=/var/log/heck_web.log
autostart=false
autorestart=true
startsecs=10
stopwaitsecs=600
